# cron/rkh.cron
#!/bin/sh
#

( /usr/bin/rkhunter --versioncheck
  /usr/bin/rkhunter --update
  /usr/bin/rkhunter --cronjob --report-warnings-only
) | /bin/mail -s 'rkhunter Daily Run' root
